# systemd/rinetd.service
[Unit]
Description=rinetd

[Service]
Type=forking
ExecStart=/usr/sbin/rinetd

[Install]
WantedBy=multi-user.target

