# systemd/rinetd.service
[Unit]
Description=rinetd

[Service]
Type=forking
Restart=on-failure
RestartSec=5s
ExecStart=/usr/sbin/rinetd

[Install]
WantedBy=multi-user.target

